# tools/youtube-downloader/bin/cron.d
0 8 * * * root /youtube-downloader/bin/run.sh
